# cuckoo-purge.conf
[program:cuckoo-purge]
command = /data/cuckoo/bin/python /data/cuckoo/.cuckoo/stuff/purge.py --cwd /data/cuckoo/.cuckoo
user = cuckoo
autostart = true
